# systemd/notes-notify.service
[Unit]
Description=Notify notes 

[Service]
ExecStart=/home/thibault/.dotfiles/scripts/notes-notify.sh

[Install]
WantedBy=default.target
